# .root/etc/profile.d/aliases.sh

alias ls="ls --color"
alias fucking=sudo
alias nv=nvim
alias copy="xclip -selection clipboard"
alias md=mkdir
alias ""="g++"
alias cargo-test="cargo test -- --nocapture"
# dotfiles
alias config='/usr/bin/git --git-dir=$HOME/.dotfiles --work-tree=$HOME'
# proxy control to set env variables
alias proxyctl=". /usr/local/bin/proxyctl"
