# .root/etc/profile.d/var.sh

# editor
export SUDO_EDITOR=/usr/bin/nvim
export EDITOR=/usr/bin/nvim
export VISUAL=/usr/bin/nvim

# miscellaneous
RUST_BACKTRACE=full
CXX=/usr/bin/clang++


